# dotastats.service
[Unit]
Description=dotastats service
Requires=docker.service
After=docker.service

[Service]
Type=simple
User=prod-user
RemainAfterExit=true
WorkingDirectory=/home/prod-user/deployment
ExecStart=/usr/bin/docker-compose up -d --remove-orphans
ExecStop=/usr/bin/docker-compose down
